# ops/v48-admin.service
# V48 admin console API (:8790) — trace/latency/failures explorer (admin/server.py). [audit R4 2026-07-12]
# Install as a user unit:
#   cp ops/v48-admin.service ~/.config/systemd/user/
#   systemctl --user daemon-reload && systemctl --user enable --now v48-admin
[Unit]
Description=V48 admin console API (:8790)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/rohith/desktop/BFI/backend/layer2/pipeline_v48
ExecStart=/home/rohith/.pyenv/versions/3.12.12/bin/python3 admin/server.py
Environment=V48_ADMIN_PORT=8790
Environment=PYTHONUNBUFFERED=1
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
